# deploy/oracle/frontend/nemorax-web.service
[Unit]
Description=Nemorax Flet web frontend
After=network-online.target nemorax-backend.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nemorax/NEMORAXS
EnvironmentFile=/etc/nemorax/web.env
ExecStart=/opt/nemorax/NEMORAXS/.venv/bin/python serve_web.py
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
